# etc/profile.d/composer-aliases.sh
alias c='composer'
alias cgr='composer global require'
alias ci='composer install'
alias cr='composer require'
alias cu='composer update'
