# etc/profile.d/composer-aliases.sh
# https://getcomposer.org/doc/03-cli.md
alias c='composer'
alias ccc='composer clear-cache'
alias ccp='composer create-project'
alias cd='composer depends'
alias cda='composer dump-autoload'
alias cg='composer global'
alias cgr='composer global require'
alias ch='composer help'
alias ci='composer install'
alias cod='composer outdated'
alias cr='composer require'
alias crd='composer require --dev'
alias crm='composer remove'
alias crmd='composer remove --dev'
alias crs='composer run-script'
alias cst='composer status'
alias cv='composer validate'
alias cu='composer update'
alias csu='composer self-update'

alias u='./vendor/bin/phpunit'
